# systemd/portage-sync.service
[Unit]
Description=Portage syncer

[Service]
Type=oneshot
ExecStart=/usr/bin/esync

[Install]
WantedBy=multi-user.target
